# etc/supervisor-conf.d/supervisord-serf-join.conf
[program:serf-join]
command=serf-join.sh
exitcodes=0
autorestart=unexpected
autostart=true
